# chat-old/powerball_project_v2/cron/powerball.cron.d
# /etc/cron.d/powerball
0 5 * * 2,4,7 root /usr/local/bin/get-powerball.sh
